# miscellaneous/devenv/celeryone.conf
description	"IK Cloud Django Development Server"

start on runlevel [2345]
stop on runlevel [!2345]

respawn
respawn limit 30 30
setuid cloud
chdir /opt/webadmin/cloud/miscellaneous/celery
env DJANGO_SETTINGS_MODULE=cloud.settings.dev                               
env DJANGO_DB_PASSWORD=asjklddfjklqjf                                       
env DJANGO_SECRET_KEY=asjklddfjklqjfasjklddfjklqjfasjklddfjklqjf
env DJANGO_SETTINGS_MODULE="" 
exec celery -A opennebula_celery worker --loglevel=info -c 1 -Q opennebula
